# debian/prerm
#!/bin/sh

invoke-rc.d wb-mqtt-smartbus stop
